# Blackberry/fast_2g_upgrade/fast-wifi.service
[Unit]
Description=Fast WiFi tune and monitor
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/tune-network.sh

[Install]
WantedBy=multi-user.target
